# polymarket-bot.service
[Unit]
Description=Polymarket Dual Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/polytest
ExecStart=/home/ubuntu/polytest/.venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
EnvironmentFile=/home/ubuntu/polytest/.env
StartLimitIntervalSec=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
